# apps/agent/deploy/coding-agent.timer
[Unit]
Description=Run Coding Agent every 15 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=15min
Persistent=true

[Install]
WantedBy=timers.target
